# deploy/systemd/agentos-backend.service
[Unit]
Description=AgentOS FastAPI Backend
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=k1212gh
WorkingDirectory=/home/k1212gh/linux-web-os/backend
EnvironmentFile=-/etc/agentos/env
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ReadWritePaths=/home/k1212gh/linux-web-os

[Install]
WantedBy=multi-user.target
